# ps3joy.conf
start on runlevel [2345]

respawn

exec /opt/ros/groovy/lib/ps3joy/ps3joy.py --inactivity-timeout=900 

